# systemd/dnf-p2p-proxy.service
[Unit]
Description=DNF P2P Package Sharing Proxy
Documentation=https://github.com/staernid/dnf-plugin-p2p
After=network.target

[Service]
Type=simple
User=dnf-p2p
Group=dnf-p2p
ExecStart=/usr/bin/python3 /usr/libexec/dnf-plugin-p2p/p2p_server.py

# Automatically restart if it crashes; not on clean exit
Restart=on-failure
RestartSec=5s

# Logging goes to journald — use: journalctl -u dnf-p2p-proxy.service
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dnf-p2p-proxy

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
CacheDirectory=dnf-plugin-p2p
ReadWritePaths=/var/cache/libdnf5 /var/cache/dnf-plugin-p2p /tmp

[Install]
WantedBy=multi-user.target
