# etc/profile.d/libdnf-python-plugin.sh
#!/bin/sh
# /etc/profile.d/libdnf-python-plugin.sh
# Enable libdnf5 Python plugin loader for dnf-plugin-p2p
export LIBDNF_PYTHON_PLUGIN_DIR=/usr/lib/python3.14/site-packages/libdnf_plugins
