# etc/profile.d/libdnf-python-plugin.sh.in
#!/bin/sh
# /etc/profile.d/libdnf-python-plugin.sh
# Enable libdnf5 Python plugin loader for dnf-plugin-p2p
export LIBDNF_PYTHON_PLUGIN_DIR=@PYTHON_INSTALL_DIR@/libdnf_plugins
